# scripts/launchd/com.yuanta.trading.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yuanta.trading.backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /path/to/local/Desktop/Project/Stock_OneClick_Order &amp;&amp; exec .venv/bin/python -m uvicorn backend.app.main:app --host 127.0.0.1 --port 8000</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/path/to/local/Desktop/Project/Stock_OneClick_Order/logs/uvicorn.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/local/Desktop/Project/Stock_OneClick_Order/logs/uvicorn.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
